# supervisord.conf
[program:frontend]
command=npm run start
autorestart=true

[program:backend]
command=/code/backend all
autorestart=true
